# systemd/.config/systemd/user/waybar-watcher.path
[Path]
PathModified=%h/.config/waybar

[Install]
WantedBy=graphical-session.target
